# config-template/supervisor.d/addons-example.conf
[program:addons-example-heartbeat]
command=/opt/openclaw-addons/bin/addons-example-heartbeat
directory=/workspace
user=appuser
autostart=true
autorestart=true
startsecs=1
priority=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
