# config-template/supervisord.conf
[unix_http_server]
file=/home/appuser/.cache/supervisor/supervisor.sock
chmod=0700
chown=appuser:appuser

[supervisord]
user=root
nodaemon=true
logfile=/home/appuser/.cache/supervisor/supervisord.log
logfile_maxbytes=0
pidfile=/home/appuser/.cache/supervisor/supervisord.pid
childlogdir=/home/appuser/.cache/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/appuser/.cache/supervisor/supervisor.sock

[program:openclaw-gateway]
command=/ops/scripts/run-openclaw-gateway.sh
directory=/app
user=appuser
autostart=true
autorestart=true
startsecs=5
priority=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:maintainer-loop]
command=/ops/scripts/run-maintainer-loop.sh
directory=/app
user=appuser
autostart=true
autorestart=true
startsecs=3
priority=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:obsidian-sync]
command=/ops/scripts/run-obsidian-sync.sh
directory=/home/appuser
user=appuser
autostart=true
autorestart=unexpected
exitcodes=0,41,42
startsecs=3
priority=30
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:obs-auto-moc-listener]
command=/ops/scripts/run-obs-auto-moc-listener.sh
directory=/home/appuser/custom-claw-tools/obs-auto-moc
user=appuser
autostart=true
autorestart=true
startsecs=2
priority=40
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:picoclaw-ops-companion]
command=/ops/scripts/run-picoclaw-ops-companion.sh
directory=/home/appuser/custom-claw-tools/picoclaw-ops-companion
user=appuser
autostart=true
autorestart=true
startsecs=2
priority=50
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:fami-ghome]
command=/ops/scripts/run-fami-ghome.sh
directory=/home/appuser/custom-claw-tools/fami-ghome
user=appuser
autostart=true
autorestart=true
startsecs=2
priority=60
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:openclaw-cron-seeder]
command=/ops/scripts/seed-openclaw-cron.sh
directory=/app
user=appuser
autostart=true
autorestart=false
startsecs=0
startretries=0
priority=70
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:supercronic]
command=/ops/scripts/run-supercronic.sh
directory=/home/appuser
user=appuser
autostart=true
autorestart=true
startsecs=0
priority=80
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
